# test-api-server.service
[Unit]
Description=Test API Server - JSON文件管理API服务器
Documentation=https://github.com/your-username/test-api-server
After=network.target
Wants=network.target

[Service]
Type=simple
User=nodejs
Group=nodejs
WorkingDirectory=/mnt/d/BaiduSyncdisk/my_code/nodejs/test_api_server
ExecStart=/usr/bin/npm start
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10

# 环境变量
Environment=NODE_ENV=production
Environment=PORT=3000

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/mnt/d/BaiduSyncdisk/my_code/nodejs/test_api_server/uploads

# 资源限制
LimitNOFILE=65536

# 日志设置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=test-api-server

[Install]
WantedBy=multi-user.target